# com.pranavpg.daily-arxiv.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pranavpg.daily-arxiv</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/pranavpg/daily-arxiv-feed/run.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/pranavpg/daily-arxiv-feed</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/pranavpg/daily-arxiv-feed/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/pranavpg/daily-arxiv-feed/logs/launchd-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/Users/pranavpg/.local/bin</string>
    </dict>
</dict>
</plist>
